# Service/jk.service
[Unit]
Description=Git Auto Sync Script

[Service]
ExecStart=/bin/bash /root/jk.sh
WorkingDirectory=/root
User=root
Restart=no
