# Service/v2ray.openrc
#!/sbin/openrc-run

# openrc service-script for V2Ray
#
# place in /etc/init.d/
# start on boot: "rc-update add v2ray"
supervisor=supervise-daemon
name="v2ray"
description="V2Ray Proxy"

command="/root/v2ray/v2ray"
command_args="-config /root/v2ray/config.json"
pidfile="/run/v2ray.pid"

command_background=true

depend() {
    need net
    after net
}